# i386/Packages/Included/frisbee_1.2-20170430-xcowsay_i386/DEBIAN/postinst
#!/bin/sh
# postinst script for frisbee
#
# see: dh_installdeb(1)

set -e

. /usr/share/debconf/confmodule

[ -f /etc/frisbee/wpa_supplicant.conf.bak ] && mv -f /etc/frisbee/wpa_supplicant.conf.bak /etc/frisbee/wpa_supplicant.conf

if [ -x /etc/init.d/dhcpcd ]; then
update-rc.d dhcpcd remove 
invoke-rc.d dhcpcd stop
fi

if [ -x "/etc/init.d/frisbee" ]; then
	update-rc.d frisbee defaults >/dev/null
	if [ -n "$2" ]; then
		_dh_action=restart
	else
		_dh_action=start
	fi
	invoke-rc.d frisbee $_dh_action || exit $?
fi



db_stop;

exit 0
