# i386/Packages/Included/frisbee_1.2-20170430-xcowsay_i386/DEBIAN/prerm
#!/bin/sh
set -e
update-rc.d -f frisbee remove
invoke-rc.d frisbee stop || exit $?




